# deploy/systemd/esotericwebb-server.service
# SPDX-License-Identifier: AGPL-3.0-or-later
# esotericwebb-server.service — esotericWebb composition server on flockGate
#
# Serves the esotericWebb interactive experience via HTTP.
# Caddy routes: webb.primals.eco → flockGate:8090.
#
# Usage: sudo systemctl enable --now esotericwebb-server
# Deploy: copy binary from plasmidBin/primals/esotericwebb
[Unit]
Description=esotericWebb Composition Server
After=songbird-federation.service membrane-nucleus.target
Wants=songbird-federation.service
ConditionPathExists=/opt/membrane/esotericwebb

[Service]
Type=simple
UMask=0002
WorkingDirectory=/opt/ecoPrimals/gardens/esotericWebb
ExecStart=/opt/membrane/esotericwebb serve --content content/ --listen 0.0.0.0:8090
Environment=ECOPRIMALS_ROOT=/opt/ecoPrimals
Environment=MEMBRANE_INSTALL_BASE=/opt/membrane
Environment=BIOMEOS_SOCKET_DIR=/run/membrane
Environment=HOME=/root
Restart=on-failure
RestartSec=3
RuntimeDirectory=membrane
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=esotericwebb

[Install]
WantedBy=membrane-nucleus.target
