# deploy/systemd/plasmid-pipeline.service
# SPDX-License-Identifier: AGPL-3.0-or-later
# plasmid-pipeline.service — Zero-touch binary refresh from depot
# Triggered by plasmid-pipeline.timer or manually via systemctl start
#
# Architecture (Wave 92): VPS does NOT build from source.
# Binaries flow: gate build → depot (plasmidBin/) → this service deploys.
# Harvest happens on builder gates (sporeGate), not here.
[Unit]
Description=Plasmid Pipeline — deploy binaries from depot to install base
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/membrane/membrane

[Service]
Type=oneshot
ExecStart=/opt/membrane/membrane plasmid.refresh
Environment=PLASMIDBIN_DEPOT=/opt/plasmidBin
Environment=MEMBRANE_INSTALL_BASE=/opt/membrane
Environment=ECOPRIMALS_ROOT=/opt/ecoPrimals
Environment=FORGEJO_SSH_HOST=git.primals.eco:2222
Environment=HOME=/root
SuccessExitStatus=1
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plasmid-pipeline
