# deploy/systemd/plasmid-pipeline.service
# plasmid-pipeline.service — Zero-touch binary harvest → refresh cycle
# Triggered by plasmid-pipeline.timer or manually via systemctl start
[Unit]
Description=Plasmid Pipeline — harvest upstream, rebuild drifted, push to VPS
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/membrane/membrane

[Service]
Type=oneshot
ExecStart=/opt/membrane/membrane plasmid.pipeline
Environment=PLASMIDBIN_DEPOT=/opt/plasmidBin
Environment=MEMBRANE_INSTALL_BASE=/opt/membrane
Environment=ECOPRIMALS_ROOT=/opt/ecoPrimals
Environment=HOME=/root
SuccessExitStatus=1
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plasmid-pipeline
