# deploy/systemd/user/membrane-nucleus@.service
# SPDX-License-Identifier: AGPL-3.0-or-later
# membrane-nucleus@.service (user) — NUCLEUS primal for desktop gates
# Instance name (%i) is the primal binary name.
# Runs under the user session, sockets in $XDG_RUNTIME_DIR/biomeos/.
#
# Before installing, set your workspace root and target triple:
#   systemctl --user set-environment ECOPRIMALS_ROOT=$HOME/Development/ecoPrimals
#   systemctl --user set-environment MEMBRANE_TARGET=x86_64-unknown-linux-musl
#
# Install: cp membrane-nucleus@.service ~/.config/systemd/user/
#          systemctl --user daemon-reload
#          systemctl --user enable --now membrane-nucleus@beardog
[Unit]
Description=NUCLEUS Primal — %i (user session)
After=songbird-federation.service
Wants=songbird-federation.service

[Service]
Type=simple
UMask=0002
ExecStart=${ECOPRIMALS_ROOT}/infra/plasmidBin/primals/${MEMBRANE_TARGET}/%i server \
    --socket %t/biomeos/%i.sock
Environment=ECOPRIMALS_ROOT=%h/Development/ecoPrimals
Environment=MEMBRANE_TARGET=x86_64-unknown-linux-musl
Environment=BIOMEOS_SOCKET_DIR=%t/biomeos
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=membrane-%i

[Install]
WantedBy=default.target
